# install/etc/s6/services/10-insync/finish
#!/bin/bash

/usr/bin/insync-headless quit
