# install/etc/s6/services/10-insync/finish
#!/bin/bash

sudo -u insync exec /usr/bin/insync-headless quit
